# One.conf
description "Yesod Masters Test"
start on runlevel [2345];
stop on runlevel [!2345];
respawn
chdir /home/ubuntu/One
exec /home/ubuntu/One/dist/build/One/One
